# etc/xdg/autostart/ubiquity-mint.desktop
[Desktop Entry]
Name=Live configuration
Comment=Live configuration
Icon=ubiquity
Exec=ubiquity-mint
Terminal=false
Type=Application
NoDisplay=true
